# cups-brother-mfc260c/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288

bashio::log.info "Warte auf CUPS-Socket..."

until [ -e /run/cups/cups.sock ]; do
    sleep 1
done

bashio::log.info "nginx-Ingress-Proxy wird gestartet"

exec nginx -g "daemon off;"
