# cups-brother-mfc260c/rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288

bashio::log.info "Starte D-Bus Daemon"

exec dbus-daemon --system --nofork
