# cups-brother-mfc260c/rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288

bashio::log.info "Warte auf D-Bus-Socket..."

until [ -e /run/dbus/system_bus_socket ]; do
    sleep 1
done

bashio::log.info "Avahi-Daemon wird gestartet"

exec avahi-daemon --no-drop-root
