# cups-brother-mfc260c/rootfs/etc/s6-overlay/s6-rc.d/cups-server/run
#!/usr/bin/with-contenv bashio

ulimit -n 524288

# ---------------------------------------------------------------------------
# USB-Quirks für Brother MFC-260C setzen
# ---------------------------------------------------------------------------
set_usb_quirks() {
    local quirk_file="/usr/share/cups/usb/org.cups.usb-quirks"

    if [ -f "${quirk_file}" ]; then
        if ! grep -qF "0x04f9 0x01d6 no-reattach" "${quirk_file}"; then
            bashio::log.info "Füge Brother MFC-260C no-reattach Quirk hinzu"
            echo "0x04f9 0x01d6 no-reattach" >> "${quirk_file}"
        fi
        if ! grep -qF "0x04f9 0x01d6 usbfs" "${quirk_file}"; then
            bashio::log.info "Füge Brother MFC-260C usbfs Quirk hinzu"
            echo "0x04f9 0x01d6 usbfs" >> "${quirk_file}"
        fi
    fi
}

# ---------------------------------------------------------------------------
# Auf Avahi-Socket warten
# ---------------------------------------------------------------------------
wait_for_avahi() {
    local timeout=30
    local elapsed=0

    bashio::log.info "Warte auf Avahi Socket..."
    until [ -e /var/run/avahi-daemon/socket ]; do
        sleep 1
        elapsed=$((elapsed + 1))
        if [ "${elapsed}" -ge "${timeout}" ]; then
            bashio::log.warning "Avahi Socket nach ${timeout}s nicht gefunden – starte CUPS trotzdem"
            return
        fi
    done
    bashio::log.info "Avahi Socket bereit"
}

# ---------------------------------------------------------------------------
# CUPS Konfiguration persistent machen
# ---------------------------------------------------------------------------
prepare_cups_config() {
    bashio::log.info "Bereite persistente CUPS-Konfiguration vor"

    if [ ! -d /config/cups ]; then
        bashio::log.info "Kopiere Standard CUPS-Konfiguration nach /config/cups"
        cp -r /etc/cups /config/cups
    fi

    rm -rf /etc/cups
    ln -s /config/cups /etc/cups
    bashio::log.info "CUPS Konfiguration verlinkt: /config/cups → /etc/cups"
}

# ---------------------------------------------------------------------------
# Hauptprogramm
# ---------------------------------------------------------------------------
set_usb_quirks
wait_for_avahi
mkdir -p /share/cups-pdf
prepare_cups_config

bashio::log.info "Starte CUPS Server"
exec cupsd -f
